# virtualize.conf
[program:virtualize]
command=/opt/parasoft/virtualize/9.8/virtualizecli -startVirtualizeServer
autostart=true
